# nginx-supervisor.ini
; Запустить supervisor в основном (не фоновом) режиме
[supervisord]
nodaemon=true

; Запустить бинарник nginx
; - user: запускать от имени пользователя root
; - autostart: запускать автоматически при старте supervisor
[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true

; Запустить бинарник php-fpm7.0
[program:php-fpm]
command = /usr/sbin/php-fpm7.0 -F
user = root
autostart = true


; Запустить cron
[program:nginx]
command = /usr/sbin/cron
user = root
autostart = true
